# hpfeeds/geoloc
### BEGIN INIT INFO
# Provides:          geoloc
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       Geo Loc Events
#
### END INIT INFO


# Using the lsb functions to perform the operations.
. /lib/lsb/init-functions

DESC="Hpfeeds Geo Locations"

# Process name ( For display )
NAME=geoloc

DIR=/opt/hpfeeds/examples/geoloc

# Daemon name, where is the actual executable
DAEMON="/usr/bin/python $DIR/geoloc.py"

# PID file for the daemon
PIDFILE=/var/run/$NAME.pid

case "$1" in
	start)
		echo -n "Starting $DESC: "
		start-stop-daemon --start --chdir $DIR --background --pidfile $PIDFILE --make-pidfile --exec $DAEMON
		echo "$NAME.pid"
	;;

	stop)
		echo -n "Stopping $DESC: "
		start-stop-daemon --stop --pidfile $PIDFILE
		echo "$NAME."
	;;

	restart)
		echo -n "Restarting $DESC: "
		start-stop-daemon --stop --pidfile $PIDFILE
		sleep 1
		start-stop-daemon --start --chdir $DIR --background --pidfile $PIDFILE --make-pidfile --exec $DAEMON
		echo "$NAME."
	;;

	status)
		if [ -e $PIDFILE ]; then
			status_of_proc -p $PIDFILE $DAEMON "$NAME process" && exit 0 || exit $?
		else
			log_daemon_msg "$NAME process is not running"
			log_end_msg 0
		fi
	;;
*)

echo "Usage: /etc/init.d/glastopf {start|stop|restart|status}"
exit 1
;;

esac
exit 0
